# docker/services.d/postgres/run
#!/command/with-contenv sh
exec s6-setuidgid postgres /usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/15/main
