# docker/services.d/redis/run
#!/command/with-contenv sh
exec redis-server --appendonly yes --dir /data --protected-mode no
